# arec.service
[Unit]
Description=AREC Audio Recording Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/arec
ExecStart=/root/arec/arec.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
